# .config/systemd/user/cliphist-clean.service
[Unit]
Description=Bersihkan riwayat cliphist setiap 30 hari

[Service]
Type=oneshot
ExecStart=/usr/bin/cliphist wipe
