# .config/systemd/user/user/wallpaper-swww-timer.service
[Unit]
Description=Ganti Wallpaper Berdasarkan Jam

[Service]
ExecStart=%h/.config/my-plugins/swww-change-by-time.sh
